# crontab
* * * * * cd /app; python ShadowbanAlerts.py
