# crontab
*/1 * * * * cd /app; pipenv run python ShadowbanAlerts.py
